# dyndns/etc/cron.d/freedns
3,8,13,18,23,28,33,38,43,48,53,58 * * * * root /usr/bin/wget --no-check-certificate -q -O - https://freedns.afraid.org/dynamic/update.php?XXX >> /var/log/dyndns/freedns.log 2>&1
